# templates/90-ethringbuffer.rules.j2
ACTION=="add", SUBSYSTEM=="net", NAME=="enp2s0", RUN+="/usr/sbin/ethtool -G $name rx 511"