# templates/heartbeat.timer.j2
[Unit]
Description=Heartbeat

[Timer]
OnActiveSec=0
OnUnitActiveSec=14m

[Install]
WantedBy=timers.target
