# templates/heartbeat.timer.j2
[Unit]
Description=Heartbeat

[Timer]
OnCalendar=*:0/12
Persistent=true

[Install]
WantedBy=timers.target
